# scripts/launchd/com.mu.hippo.ingest-v2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mu.hippo.ingest-v2</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/mu/src/hippo-public/scripts/ingest_v2.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mu/src/hippo-public</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mu/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mu</string>
        <key>HIPPO_HOME</key>
        <string>/Users/mu/.hippo</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mu/.hippo/logs/ingest.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mu/.hippo/logs/ingest.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
